# docker-alpine/etc/init.d/c8y-configuration-plugin
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Thin-edge device configuration management for Cumulocity"

command="/bin/c8y_configuration_plugin"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
output_log="/var/log/$RC_SVCNAME/$RC_SVCNAME.log"
error_log="/var/log/$RC_SVCNAME/$RC_SVCNAME.err"

depend() {
    need net
}

start_pre() {
    /bin/mkdir -m 744 -p /var/log/$RC_SVCNAME
}

